# tools/diagnostic/cron/runmystore-diagnostic.crontab
# /etc/cron.d/runmystore-diagnostic
# RunMyStore.AI Diagnostic Framework
TZ=Europe/Sofia
SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

# Nightly sales pulse (5-15 random sales на tenant=7)
0  3 * * *   www-data /var/www/runmystore/tools/diagnostic/cron/sales_pulse.sh

# Weekly diagnostic — понеделник 03:00
0  3 * * 1   www-data /var/www/runmystore/tools/diagnostic/cron/diagnostic_weekly.sh

# Monthly diagnostic — 1-ви ден 04:00
0  4 1 * *   www-data /var/www/runmystore/tools/diagnostic/cron/diagnostic_monthly.sh

# Daily summary email — 08:30 ВСЕКИ ден
30 8 * * *   www-data /var/www/runmystore/tools/diagnostic/cron/daily_summary.sh
